# files/supervisord/supervisord_env.conf
[program:supervisor-proc-exit-listener]
command=/usr/bin/supervisor-proc-exit-listener --container-name=%(ENV_CONTAINER_NAME)s
autorestart=unexpected
exitcodes=0
startsecs=0
events=PROCESS_STATE
